# ncdu/.config/systemd/user/ncdu.service
[Unit]
Description=NCDU scanner

[Service]
ExecStart=/usr/bin/ncdu -o %C/ncdu

[Install]
WantedBy=default.target
