# ncdu/.config/systemd/user/ncdu.timer
[Unit]
Description=Run NCDU every six hours

[Timer]
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
